# nas/etc/sv/plex/run
#!/bin/sh

exec /usr/bin/podman run \
  --replace \
  --rm \
  --name=plex \
  -p 32400:32400 \
  --device /dev/dri:/dev/dri \
  -v /data/appdata/plex:/config:Z \
  -v /data/media:/media:z \
  ocker.io/plexinc/pms-docker:latest