# nas/etc/sv/komga/run
#!/bin/sh

exec /usr/bin/podman run \
  --replace \
  --rm \
  --name=komga \
  -p 25600:25600 \
  -v /data/appdata/komga:/config:Z \
  -v /data/books:/data:ro,z \
  docker.io/gotson/komga:latest