# nas/huai/.config/systemd/user/rclone-aliyun.service
[Unit]
# Optimized for Aliyun/WebDAV 1080p streaming
Description=rclone aliyun User Mount Read-Only
After=network-online.target

[Service]
Type=simple
# Emby requires allow-other to read fuse mounts from podman
ExecStart=/usr/bin/rclone mount aliyun: /data/media/aliyun \
    --config=/home/huai/.config/rclone/rclone.conf \
    --read-only \
    --allow-other \
    --cache-dir /podman/rclone-cache \
    --vfs-cache-mode full \
    --vfs-cache-max-size 10G \
    --vfs-cache-max-age 24h \
    --dir-cache-time 1000h \
    --attr-timeout 1000h \
    --poll-interval 0 \
    --buffer-size 128M \
    --vfs-read-ahead 256M \
    --vfs-read-chunk-size 64M \
    --vfs-read-chunk-size-limit 1G \
    --no-modtime \
    --tpslimit 10 \
    --log-level ERROR
ExecStop=/usr/bin/fusermount -u /data/media/aliyun
# Use none to let fusermount handle the unmount properly
KillMode=none
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target