# nas/etc/sv/pw/run
#!/bin/sh

exec /usr/bin/podman run \
  --replace \
  --rm \
  --init \
  --name=pw \
  -p 8181:80 \
  -v /data/appdata/vaultwarden:/data:Z \
  docker.io/vaultwarden/server:latest