# openwrt/openwrt/etc/init.d/tproxy
#!/bin/sh /etc/rc.common

START=99
STOP=10
USE_PROCD=1

# --- Configuration ---
PROXY_PORT=7893
ROUTING_MARK=666
PROXY_FWMARK=0x1
TABLE_ID=100
# --- End Configuration ---

boot() {
    ( sleep 10 && start_service ) &
}

service_triggers() {
    procd_add_reload_trigger "firewall"
    procd_add_reload_trigger "network"
}

start_service() {
    # 使用 logger 记录开始尝试加载
    logger -t singbox-tproxy "Starting singbox-tproxy rules..."

    procd_open_instance
    procd_set_param command /bin/sh -c '
        NFT_TABLE="inet fw4"
        CHAIN_PRE="singbox_prerouting"
        CHAIN_OUT="singbox_output"
        
        # 1. Create Chains
        nft list chain $NFT_TABLE $CHAIN_PRE >/dev/null 2>&1 || nft add chain $NFT_TABLE $CHAIN_PRE
        nft flush chain $NFT_TABLE $CHAIN_PRE
        
        nft list chain $NFT_TABLE $CHAIN_OUT >/dev/null 2>&1 || nft add chain $NFT_TABLE $CHAIN_OUT
        nft flush chain $NFT_TABLE $CHAIN_OUT

        # 2. PREROUTING Rules (Safe comments)
        nft add rule $NFT_TABLE $CHAIN_PRE meta mark '$ROUTING_MARK' return comment "Prevention_Bypass_Singbox_Outbound"
        nft add rule $NFT_TABLE $CHAIN_PRE udp dport { 67, 68 } return comment "Bypass_DHCP"
        nft add rule $NFT_TABLE $CHAIN_PRE ip6 nexthdr icmpv6 return comment "Bypass_IPv6_ICMP"
        nft add rule $NFT_TABLE $CHAIN_PRE ct status dnat return comment "Bypass_Port_Forwarding_DNAT"
        nft add rule $NFT_TABLE $CHAIN_PRE fib daddr type local return comment "Bypass_Router_Local_Access"
        nft add rule $NFT_TABLE $CHAIN_PRE ip daddr { 127.0.0.0/8, 10.0.0.0/8, 172.16.0.0/12, 192.168.0.0/16, 224.0.0.0/4 } return
        nft add rule $NFT_TABLE $CHAIN_PRE ip6 daddr { ::1/128, fc00::/7, fe80::/10, ff00::/8 } return
        nft add rule $NFT_TABLE $CHAIN_PRE meta l4proto { tcp, udp } tproxy to :'$PROXY_PORT' meta mark set '$PROXY_FWMARK' accept comment "Action_TProxy_To_Singbox"

        # 3. OUTPUT Rules
        nft add rule $NFT_TABLE $CHAIN_OUT meta oifname "lo" return comment "Bypass_Loopback_Interface"
        nft add rule $NFT_TABLE $CHAIN_OUT meta mark '$ROUTING_MARK' return comment "Prevention_Bypass_Singbox_Outbound"
        nft add rule $NFT_TABLE $CHAIN_OUT ip6 nexthdr icmpv6 return comment "Bypass_IPv6_ICMP"
        nft add rule $NFT_TABLE $CHAIN_OUT udp dport { 137, 138 } return comment "Bypass_NetBIOS"
        nft add rule $NFT_TABLE $CHAIN_OUT ip daddr { 127.0.0.0/8, 10.0.0.0/8, 172.16.0.0/12, 192.168.0.0/16, 224.0.0.0/4 } return
        nft add rule $NFT_TABLE $CHAIN_OUT ip6 daddr { ::1/128, fc00::/7, fe80::/10, ff00::/8 } return
        nft add rule $NFT_TABLE $CHAIN_OUT meta l4proto { tcp, udp } meta mark set '$PROXY_FWMARK' accept comment "Action_Mark_Router_Output"

        # 4. Insert Hooks
        nft list chain $NFT_TABLE mangle_prerouting | grep -q "jump $CHAIN_PRE" || \
            nft insert rule $NFT_TABLE mangle_prerouting jump $CHAIN_PRE comment "Singbox_Pre_Hook"
        nft list chain $NFT_TABLE mangle_output | grep -q "jump $CHAIN_OUT" || \
            nft insert rule $NFT_TABLE mangle_output jump $CHAIN_OUT comment "Singbox_Out_Hook"

        # 5. Routing
        ip route show table '$TABLE_ID' | grep -q "local default" || ip route add local default dev lo table '$TABLE_ID'
        ip rule show table '$TABLE_ID' | grep -q "fwmark '$PROXY_FWMARK'" || ip rule add fwmark '$PROXY_FWMARK' lookup '$TABLE_ID' pref 9000

        if [ -f /proc/net/if_inet6 ]; then
            ip -6 route show table '$TABLE_ID' | grep -q "local default" || ip -6 route add local default dev lo table '$TABLE_ID'
            ip -6 rule show table '$TABLE_ID' | grep -q "fwmark '$PROXY_FWMARK'" || ip -6 rule add fwmark '$PROXY_FWMARK' lookup '$TABLE_ID' pref 9000
        fi
        
        # 成功后写入日志
        logger -t tproxy "Rules loaded successfully (Safe Mode)"
    '
    procd_set_param stdout 1
    procd_set_param stderr 1
    procd_close_instance
}

stop_service() {
    NFT_TABLE="inet fw4"
    CHAIN_PRE="singbox_prerouting"
    CHAIN_OUT="singbox_output"

    # 1. Remove Hooks
    for h in $(nft --handle list chain $NFT_TABLE mangle_prerouting | grep "Singbox_Pre_Hook" | awk "{print \$NF}"); do
        nft delete rule $NFT_TABLE mangle_prerouting handle $h 2>/dev/null
    done
    for h in $(nft --handle list chain $NFT_TABLE mangle_output | grep "Singbox_Out_Hook" | awk "{print \$NF}"); do
        nft delete rule $NFT_TABLE mangle_output handle $h 2>/dev/null
    done

    # 2. Flush chains
    nft flush chain $NFT_TABLE $CHAIN_PRE 2>/dev/null; nft delete chain $NFT_TABLE $CHAIN_PRE 2>/dev/null
    nft flush chain $NFT_TABLE $CHAIN_OUT 2>/dev/null; nft delete chain $NFT_TABLE $CHAIN_OUT 2>/dev/null

    # 3. Clear Routing
    ip route flush table $TABLE_ID 2>/dev/null
    ip rule del fwmark $PROXY_FWMARK table $TABLE_ID 2>/dev/null
    [ -f /proc/net/if_inet6 ] && { ip -6 route flush table $TABLE_ID 2>/dev/null; ip -6 rule del fwmark $PROXY_FWMARK table $TABLE_ID 2>/dev/null; }
    
    # 这里使用 logger 确保一定会写入系统日志
    logger -t tproxy "Rules and routes cleared successfully"
}